# pycheribuild/files/minimal-image/need-cplusplus.files
# These files depend on C++ support which may not be available:
etc/rc.d/devd
etc/devd.conf
